# bot.service.example
[Unit]
Description=Telegram Educational Bot
After=network.target

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/opt/telegram-bot
EnvironmentFile=/opt/telegram-bot/.env
ExecStart=/opt/telegram-bot/.venv/bin/python /opt/telegram-bot/run_bot.py
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target